# launchd/com.telegram_agent.worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.telegram_agent.worker</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/server/ai_projects/telegram_agent/.venv/bin/python3</string>
        <string>/Users/server/ai_projects/telegram_agent/worker_queue.py</string>
        <string>--interval</string>
        <string>10</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/server/ai_projects/telegram_agent</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/server/ai_projects/telegram_agent/logs/worker.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/server/ai_projects/telegram_agent/logs/worker_error.log</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
